# aws_terraform/environments/Dev/ansible/roles/webservers/files/consumer.service
[Unit]
Description=Kafka Consumer Service

[Service]
ExecStart=/home/ec2-user/consumer
WorkingDirectory=/home/ec2-user/