# aws_terraform/environments/Dev/ansible/roles/zookeeper/files/producer.service
[Unit]
Description=Kafka Producer Service

[Service]
ExecStart=/home/ec2-user/producer
WorkingDirectory=/home/ec2-user/